# projects/sandbox-sync-agent/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

# Sandbox 主进程（code-server）
# 通过环境变量 FASTGPT_ENABLE_CODE_SERVER=true|false 控制是否启动，默认 true
[program:code-server]
command=/home/sandbox/entrypoint.sh
user=sandbox
autostart=%(ENV_FASTGPT_ENABLE_CODE_SERVER)s
autorestart=true
stdout_logfile=/var/log/supervisor/code-server.log
stderr_logfile=/var/log/supervisor/code-server-error.log

# Sync Agent 后台进程（始终启动）
[program:sync-agent]
command=/opt/sync-agent/sync.sh
user=sandbox
autostart=true
autorestart=true
environment=FASTGPT_SYNC_PATH="%(ENV_FASTGPT_WORKDIR)s",HTTP_SERVER_PATH="/opt/sync-agent/http_server.py"
stdout_logfile=/var/log/supervisor/sync-agent.log
stderr_logfile=/var/log/supervisor/sync-agent-error.log
